# scripts/mathrankapi.service
[Unit]
Description=Mathrank API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/api
ExecStart=/usr/bin/java -jar /home/ubuntu/api/app.jar
Restart=on-failure
RestartSec=5
SuccessExitStatus=143
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
